# task2/local/crontab/README
# /etc/crontab

07 7 * * * /opt/php56/bin/php -f /var/www/test_ratio/data/www/ratio-test.asperans.ru/local/crontab/deffered_items.php >
/var/www/test_ratio/data/www/ratio-test.asperans.ru/local/crontab/logs/deffered_items.log